# deploy/systemd/cloudflare-tunnel.service
[Unit]
Description=Cloudflare Tunnel — Job Pipeline API
Documentation=https://github.com/rutvijmavani/auto-email
After=network-online.target pipeline-api.service
Wants=network-online.target
Requires=pipeline-api.service
StartLimitBurst=5
StartLimitIntervalSec=300s

[Service]
Type=simple
User=opc
Group=opc
WorkingDirectory=/home/opc/mail
EnvironmentFile=/home/opc/mail/.env
ExecStart=/home/opc/mail/venv/bin/python scripts/tunnel_manager.py
Restart=always
RestartSec=15s
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudflare-tunnel

[Install]
WantedBy=multi-user.target
