# deploy/systemd/recruiter-watchdog.service
[Unit]
Description=Recruiter Pipeline Watchdog (health monitoring + alerts + self-healing)
Documentation=https://github.com/rutvijmavani/auto-email
# See recruiter-scheduler.service for the network-online.target prerequisite note.
# Same requirement applies here: verify systemd-networkd-wait-online or
# NetworkManager-wait-online is enabled before relying on this ordering.
After=network-online.target
Wants=network-online.target

# Watchdog is intentionally NOT After=recruiter-scheduler.service so it keeps
# running and alerting even when the scheduler is down.

# Soft-require Redis and Postgres — if they're not systemd-managed (e.g. managed
# by the cloud provider), remove the relevant Wants/After lines.
Wants=redis.service
After=redis.service
Wants=postgresql.service
After=postgresql.service

# ── Escalation — repeated rapid failures ─────────────────────────────────────
# StartLimitBurst/StartLimitIntervalSec belong in [Unit] (systemd v229+).
StartLimitBurst=5
StartLimitIntervalSec=300s
OnFailure=recruiter-pipeline-alert@%n.service

[Service]
Type=simple
User=opc
Group=opc
WorkingDirectory=/home/opc/mail
EnvironmentFile=/home/opc/mail/.env

ExecStart=/home/opc/mail/venv/bin/python -m workers.watchdog

# ── Restart policy ────────────────────────────────────────────────────────────
# Faster restart than scheduler — watchdog is lightweight and should
# always be running to catch scheduler failures.
Restart=always
RestartSec=10s

# ── Sandboxing ───────────────────────────────────────────────────────────────
# NoNewPrivileges is intentionally omitted: the watchdog calls
# "sudo systemctl restart recruiter-scheduler" for self-healing, which requires
# the setuid bit on /usr/bin/sudo to work. NoNewPrivileges=true disables setuid
# and would break that path.
PrivateTmp=true
# read-only (not strict) — WorkingDirectory is under /home
ProtectHome=read-only

# ── Logging ──────────────────────────────────────────────────────────────────
StandardOutput=journal
StandardError=journal
SyslogIdentifier=recruiter-watchdog

[Install]
WantedBy=multi-user.target
